# sv/mongod/run
#!/bin/sh
exec 2>&1
exec sudo -u mongodb /usr/bin/mongod -f /etc/mongod.conf
